# .config/autostart/compiz.desktop

[Desktop Entry]
Type=Application
Name=Compiz
Exec=compiz --replace --indirect-rendering --loose-binding ccp
Icon=system-run
Comment=
