# .config/autostart/gtk-window-decorator.desktop

[Desktop Entry]
Type=Application
Name=GTK Window Decorator
Exec=gtk-window-decorator
Icon=system-run
Comment=
